# scripts/com.argus.rebuild-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.argus.rebuild-watcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/burtonstuff/Argus1/scripts/argus-rebuild-watcher.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>30</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/burtonstuff/.argus/rebuild-watcher.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/burtonstuff/.argus/rebuild-watcher.log</string>
</dict>
</plist>
